# service/riak/run
#!/bin/sh
exec /usr/local/riak/bin/riak console
